# deadeye/conf/deadeye.timer
[Unit]
Description=Deadeye Start Timer

[Timer]
OnStartupSec=15

[Install]
WantedBy=timers.target
