# conf/deadeye.service
[Unit]
Description=Deadeye Vision System

[Service]
ExecStart=/usr/local/sbin/deadeye
Type=simple

[Install]
WantedBy=default.target  
